# pitakbot/bot.conf
[program:bot]
directory=/home/pitak/pitakproject
command=python3 manage.py runbot

autostart=true
autorestart=true

stderr_logfile=/home/pitak/logfile_err.log
stdout_logfile=/home/pitak/logfile.log
